# misc/crontab
0 4 * * * /home/www/ccis/research/prl/generator/cron-buid.sh
